# systemd/emitter.service.in
[Unit]
Description=EPICS Relay UDP Emitter
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=20
User=root
EnvironmentFile=@CMAKE_INSTALL_PREFIX@/etc/epics-relay.conf
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/emitter \
            ${EPICS_RELAY_EMITTER_OPTIONS} \
            ${EPICS_RELAY_EMITTER_IFACE}

[Install]
WantedBy=multi-user.target
