# systemd/epics_udp_emitter@.service.in
[Unit]
Description=EPICS Relay UDP Emitter (%i)
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=20
User=root
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/epics_udp_emitter -c @CMAKE_INSTALL_FULL_SYSCONFDIR@/epics-relay_%i.conf

[Install]
WantedBy=multi-user.target
