# systemd/epics_udp_collector.service.in
[Unit]
Description=EPICS Relay UDP Collector
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=20
User=root
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/epics_udp_collector

[Install]
WantedBy=multi-user.target
